# launchd/com.krishna.agentmemory.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.krishna.agentmemory.snapshot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>scripts/mechanical_snapshot.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/krishna.champaneria/Desktop/agent-memory</string>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/agentmemory.snapshot.err</string>
  <key>StandardOutPath</key>
  <string>/tmp/agentmemory.snapshot.out</string>
</dict>
</plist>
